# rpi/resources/wifi-watchdog/wifi-watchdog.service
[Unit]
Description=WiFi USB Watchdog Service
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash /opt/wifi-watchdog/wifi-watchdog.sh
Restart=always
RestartSec=10
SyslogIdentifier=wifi-watchdog
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
